# user/anki.service
[Unit]
Description=Anki
After=xorg.target

[Service]
ExecStart=/usr/bin/anki
KillMode=process