# user/app@.service
[Unit]
Description=Application "%i"
After=xorg.target
IgnoreOnIsolate=yes

[Service]
ExecStart=/usr/bin/env %i

KillMode=process