# user/emacsclient@.service
[Unit]
Description=Emacsclient %I
After=xorg.target

[Service]
ExecStart=/usr/bin/emacsclient -c
