# user/geary.service
[Unit]
Description=Geary Mail Client
After=xorg.target

[Service]
ExecStart=/usr/bin/geary
KillMode=process