# user/gpa.service
[Unit]
Description=GPA
After=xorg.target

[Service]
ExecStart=/usr/bin/gpa
KillMode=process