# user/pidgin.service
[Unit]
Description=Pidgin
After=xorg.target

[Service]
ExecStart=/usr/bin/pidgin
KillMode=process