# user/quodlibet.service
[Unit]
Description=Quodlibet
After=xorg.target

[Service]
ExecStart=/usr/bin/quodlibet
KillMode=process