# user/seahorse.service
[Unit]
Description=Seahorse
After=xorg.target

[Service]
ExecStart=/usr/bin/seahorse
KillMode=process
