# user/skype.service
[Unit]
Description=Skpe	
After=xorg.target

[Service]
ExecStart=/opt/bin/skype
KillMode=process