# user/transmission.service
[Unit]
Description=Transmission
After=xorg.target

[Service]
ExecStart=/usr/bin/transmission-gtk
KillMode=process