# user/udiskie.service
[Unit]
Description=Udiskie automounter
After=xorg.target

[Service]
ExecStart=/usr/bin/udiskie

[Install]
WantedBy=default.target
