# user/urxvt@.service
[Unit]
Description=Urxvt %I
After=xorg.target

[Service]
ExecStart=/usr/bin/urxvt
