# user/zathura.service
[Unit]
Description=Zathura
After=xorg.target

[Service]
ExecStart=/usr/bin/zathura
KillMode=process